# scripts/cifsd.service
[Unit]
Description=CIFSD userspace service

[Service]
Type=forking
ExecStart=/usr/local/sbin/cifsd -s

[Install]
WantedBy=multi-user.target
